# lesson5/part1/testlog.timer
[Unit]
Description=Testlog timer
Requires=testlog.service

[Timer]
Unit=testlog.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target